# installer/resources/templates/agent-player-frontend.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.agentplayer.frontend</string>

    <key>ProgramArguments</key>
    <array>
        <string>{{INSTALL_DIR}}/runtime/node/bin/node</string>
        <string>{{INSTALL_DIR}}/app/packages/frontend/server.js</string>
    </array>

    <key>WorkingDirectory</key>
    <string>{{INSTALL_DIR}}/app</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>NODE_ENV</key>
        <string>production</string>
        <key>PORT</key>
        <string>{{FRONTEND_PORT}}</string>
        <key>NEXT_PUBLIC_API_URL</key>
        <string>http://localhost:{{BACKEND_PORT}}</string>
        <key>PATH</key>
        <string>{{INSTALL_DIR}}/runtime/node/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>OtherJobEnabled</key>
        <dict>
            <key>com.agentplayer.backend</key>
            <true/>
        </dict>
    </dict>

    <key>StandardOutPath</key>
    <string>{{INSTALL_DIR}}/logs/frontend.log</string>

    <key>StandardErrorPath</key>
    <string>{{INSTALL_DIR}}/logs/frontend-error.log</string>

    <key>UserName</key>
    <string>{{USER}}</string>

    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
