# installer/resources/templates/agent-player.service.template
[Unit]
Description=Agent Player Backend Service
After=network.target
Wants=agent-player-frontend.service

[Service]
Type=simple
User={{USER}}
WorkingDirectory={{INSTALL_DIR}}/app/packages/backend
Environment="NODE_ENV=production"
Environment="PORT={{BACKEND_PORT}}"
Environment="PATH={{INSTALL_DIR}}/runtime/node/bin:{{INSTALL_DIR}}/runtime/python/bin:$PATH"
Environment="PYTHON_PATH={{INSTALL_DIR}}/runtime/python/bin/python3"
Environment="DATABASE_PATH={{INSTALL_DIR}}/.data/agent-player.db"
Environment="STORAGE_BASE_DIR={{INSTALL_DIR}}/public/storage"
Environment="MEMORY_BASE_DIR={{INSTALL_DIR}}/.data/memory"
ExecStart={{INSTALL_DIR}}/runtime/node/bin/node dist/api/server.js
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=multi-user.target
