# installer/resources/templates/agent-player.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.agentplayer.backend</string>

    <key>ProgramArguments</key>
    <array>
        <string>{{INSTALL_DIR}}/runtime/node/bin/node</string>
        <string>{{INSTALL_DIR}}/app/packages/backend/dist/api/server.js</string>
    </array>

    <key>WorkingDirectory</key>
    <string>{{INSTALL_DIR}}/app/packages/backend</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>NODE_ENV</key>
        <string>production</string>
        <key>PORT</key>
        <string>{{BACKEND_PORT}}</string>
        <key>PATH</key>
        <string>{{INSTALL_DIR}}/runtime/node/bin:{{INSTALL_DIR}}/runtime/python/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>StandardOutPath</key>
    <string>{{INSTALL_DIR}}/logs/backend.log</string>

    <key>StandardErrorPath</key>
    <string>{{INSTALL_DIR}}/logs/backend-error.log</string>

    <key>UserName</key>
    <string>{{USER}}</string>

    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
