# installer/resources/templates/agent-player-frontend.service.template
[Unit]
Description=Agent Player Frontend Service
After=network.target agent-player.service
Requires=agent-player.service

[Service]
Type=simple
User={{USER}}
WorkingDirectory={{INSTALL_DIR}}/app
Environment="NODE_ENV=production"
Environment="PORT={{FRONTEND_PORT}}"
Environment="NEXT_PUBLIC_API_URL=http://localhost:{{BACKEND_PORT}}"
Environment="PATH={{INSTALL_DIR}}/runtime/node/bin:$PATH"
ExecStart={{INSTALL_DIR}}/runtime/node/bin/node packages/frontend/server.js
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=multi-user.target
